# misc/shairport-sync/root/etc/services.d/shaiport/run
#!/usr/bin/with-contenv bash

exec \
shairport-sync
