# alpine-www/NGINX/rootfs/etc/services.d/nginx/run
#!/bin/execlineb -P

exec -c
nginx -g "daemon off;"
