# webapps/etherpad/v1.7/root/etc/services.d/etherpad/run
#!/usr/bin/with-contenv sh

export NODE_ENV=production
exec /opt/etherpad/bin/safeRun.sh /dev/stdout --root
