# db/redis/root/etc/services.d/redis/run
#!/bin/execlineb -P

exec -c
s6-setuidgid abc redis-server /config/redis/conf/redis.conf --protected-mode no
