# alpine-misc/jabbergram/rootfs/etc/services.d/jabbergram/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
exec -c
/usr/bin/jabbergram.py /etc/jabbergram/config.ini
