# alpine-cache/memcached/rootfs/etc/services.d/memcached-error-log/run
#!/usr/bin/execlineb -P

logutil-service -f /var/run/s6/memcached-error-log-fifo /var/log/memcached-error-logs
