# alpine-www/hiawatha/rootfs/etc/services.d/hiawatha/run
#!/bin/execlineb -P

exec -c
hiawatha -d
