# misc/prosody/root/etc/services.d/prosody/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
exec -c
/usr/bin/prosodyctl --config /config/prosody/conf/prosody.cfg.lua --root start
