# www/traefik/root/etc/services.d/traefik/run
#!/bin/execlineb -P
fdmove -c 2 1
exec -c
/usr/sbin/traefik -c /config/traefik/config/traefik.toml
