# alpine-webapps/etherpad/rootfs/etc/services.d/etherpad/run
#!/bin/execlineb -P

exec -c
 /usr/share/webapps/etherpad/bin/run.sh --root
