# alpine-cache/memcached/rootfs/etc/services.d/memcached-general-log/run
#!/usr/bin/execlineb -P

logutil-service -f /var/run/s6/memcached-general-log-fifo /var/log/memcached-general-logs
