# alpine-www/php-fpm/rootfs/etc/services.d/php-fpm/run
#!/bin/execlineb -P

exec -c
php-fpm --nodaemonize
