# alpine-db/redis/rootfs/etc/services.d/redis/run
#!/bin/execlineb -P

exec -c
s6-setuidgid redis redis-server --protected-mode no
