# alpine-misc/transmission/rootfs/etc/services.d/transmission/run
#!/bin/execlineb -P

exec -c
transmission-daemon --foreground --config-dir /var/transmission/config
