# alpine-webapps/botpress/rootfs/etc/services.d/botpress/run
#!/usr/bin/execlineb -P
cd /bot/BOTNAME
botpress start
